# scripts/com.platform.loop-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.platform.loop-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/andymontgomery/projects/platform/scripts/loop_daemon.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/andymontgomery/projects/platform</string>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/andymontgomery/Library/Logs/platform-loop-daemon.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/andymontgomery/Library/Logs/platform-loop-daemon.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/andymontgomery</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
